# openrc/init.d/run-workload
#!/sbin/openrc-run

command="/bin/run-workload-reboot"
command_background="no"
pidfile="/tmp/workload.pid"