# openrc/init.d/run-workload-stream
#!/sbin/openrc-run

command="/bin/stream"
command_background="yes"
pidfile="/tmp/stream.pid"