# openrc/init.d/run-workload-dd
#!/sbin/openrc-run

command="/bin/dd-workload.sh"
command_background="yes"
pidfile="/tmp/dd-workload.pid"